# debian/cron.d
# Launch the command to empty the spool
* * * * * root nagios2mantis empty

