# debian/cron.d
# Launch the command to empty the spool
* * * * * nagios nagios2mantis empty
