# debian/cron.d
# Launch the command to empty the spool
* * * * * nagios nagios2mantis empty
0 2 * * * nagios nagios2mantis clean
